# deploy/systemd/community_helper_index_refresh.service
[Unit]
Description=Community Helper Vector Index Refresh
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/var/www/Community_helper
Environment="PATH=/var/www/Community_helper/.venv/bin"
Environment="HF_HOME=/root/.cache/huggingface"
Environment="TRANSFORMERS_CACHE=/root/.cache/huggingface"
Environment="HF_HUB_OFFLINE=1"
Environment="TRANSFORMERS_OFFLINE=1"
ExecStart=/var/www/Community_helper/scripts/rebuild_vector_index.sh

# 索引文件更新后，retriever 会在后续请求中直接读取新文件，因此这里不需要重启 Web 服务。

